# supervisor_program.ini.erb
;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;
; PUPPET MANAGED                                          ;
; Do not edit this file on a server node unless you       ;
; are willing to have your changes overwritten by         ;
; Puppet.  If you really want to change the contents      ;
; of this file, change it in the puppet subversion        ;
; repository and check it out on the ops server.          ;
;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;
[program:<%= name -%>]
command=<%= command %>
<% if numprocs.to_i > 1 -%>
process_name=%(program_name)s_%(process_num)02d
<% end -%>
numprocs=<%= numprocs %>
<% unless priority == "" -%>
priority=<%= priority %>
<% end -%>
autostart=<%= autostart %>
autorestart=<%= autorestart %>
startsecs=<%= startsecs %>
startretries=<%= retries %>
exitcodes=<%= exitcodes %>
stopsignal=<%= stopsignal %>
stopwaitsecs=<%= stopwait %>
<% unless user == "" -%>
user=<%= user %>
<% end -%>
redirect_stderr=<%= redirect_stderr %>
stdout_logfile=<%= log_path %>/<%= name -%>.out
stdout_logfile_maxbytes=<%= stdout_logfile_maxsize %>
stdout_logfile_backups=<%= stdout_logfile_keep %>
stderr_logfile=<%= log_path %>/<%= name -%>.err
stderr_logfile_maxbytes=<%= stderr_logfile_maxsize %>
stderr_logfile_backups=<%= stderr_logfile_keep %>
<% unless environment == "" -%>
environment=<%= environment %>
<% end -%>
<% unless chdir == "" -%>
directory=<%= chdir %>
<% end -%>
<% unless umask == "" %>
umask=<%= umask %>
<% end -%>
